# branches/0.9/zumastor/debian/zumastor.postrm
#!/bin/sh
if [ "$1" = "purge" ] ; then
	rm -rf /var/lib/zumastor
	update-rc.d zumastor remove >/dev/null || exit 0
fi
