# zumastor/etc/cron.weekly/zumastor
#!/bin/bash
for vol in $(ls /var/run/zumastor/cron) ; do
	zumastor snapshot $vol weekly
done
